# supervisor/chatty-ai.ini
[program:chatty-ai]
command=sh -c "docker compose up"
process_name=chatty-ai
autostart=true
autorestart=true
directory=/root/web/chatty-ai
stdout_logfile=~/web/chatty-ai/output.log
stderr_logfile=~/web/chatty-ai/error.log